# server/systemd/altair-update.service
[Unit]
Description=altair-brain — aggiornamento periodico (fetch+reset al remoto)
After=network-online.target
Wants=network-online.target
OnFailure=altair-notify@%n.service

[Service]
Type=oneshot
User=altair
EnvironmentFile=/home/altair/altair-api/.env
ExecStart=/usr/bin/bash /home/altair/altair-api/update_brain.sh
